# download_manager/download-manager.service
[Unit]
Description=Download Manager for Agent Network
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/download_manager
ExecStart=/usr/bin/python3 download_manager.py --watch
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
